# deploy/supervisor/pdf-workers.conf
; Supervisor — chạy song song Quyển 1 + Quyển 2
; Copy vào /etc/supervisor/conf.d/ rồi: supervisorctl reread && supervisorctl update
;
; Điều chỉnh:
;   - command: đường dẫn php và thư mục project
;   - user: www-data hoặc user deploy

[program:laravel-pdf-q1]
process_name=%(program_name)s
command=php /var/www/vuongkimbao/artisan queue:work --queue=pdf-q1 --timeout=600 --memory=512 --tries=2 --sleep=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/vuongkimbao/storage/logs/pdf-worker-q1.log
stopwaitsecs=620

[program:laravel-pdf-q2]
process_name=%(program_name)s
command=php /var/www/vuongkimbao/artisan queue:work --queue=pdf-q2 --timeout=600 --memory=512 --tries=2 --sleep=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/vuongkimbao/storage/logs/pdf-worker-q2.log
stopwaitsecs=620

; Job khác (mail, default…) — tùy chọn
[program:laravel-worker-default]
process_name=%(program_name)s
command=php /var/www/vuongkimbao/artisan queue:work --queue=default --timeout=120 --memory=256 --tries=3 --sleep=3
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/vuongkimbao/storage/logs/worker-default.log
